# private_dot_config/systemd/user/tray-pre.target
[Unit]
Description = "Precursor to tray.target"
Before=tray.target

[Install]
RequiredBy=tray.target
